# autostart/systemd/fenrir.service
[Unit]
Description=Fenrir screenreader

[Service]
ExecStart=/usr/bin/fenrir
Type=forking

[Install]
WantedBy=default.target
