# autostart/systemd/fenrir.service
[Unit]
Description=fenrir
#After=

[Service]
ExecStart=path_to_executable
Type=forking
PIDFile=path_to_pidfile

[Install]
WantedBy=default.target
